# conf/systemd/agoatalarm.service.in
[Unit]
Description=ago control ATAlarm interface

[Service]
ExecStart=@BINDIR@/agoatalarm
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
