# conf/systemd/agoaxisaudiotransmit.service.in
[Unit]
Description=ago control axis audio transmit device

[Service]
ExecStart=@BINDIR@/agoaxisaudiotransmit.py
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
