# conf/systemd/agolifx.service.in
[Unit]
Description=ago control LIFX handler

[Service]
ExecStart=@BINDIR@/agolifx.py
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
