# conf/systemd/agoonvif.service.in
[Unit]
Description=ago control for onvif devices

[Service]
ExecStart=@BINDIR@/agoonvif.py
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
