# conf/systemd/agotellstick.service.in
[Unit]
Description=ago control Tellstick handler

[Service]
ExecStart=@BINDIR@/agotellstick.py
WorkingDirectory=@BINDIR@
StandardOutput=syslog
StandardError=syslog
Restart=always
RestartSec=2
User=agocontrol

[Install]
WantedBy=multi-user.target
